# root/etc/services.d/retroassembly/run
#!/bin/sh
echo "[INFO] Starting RetroAssembly..."
cd /app
exec s6-setuidgid bsd node dist/server/serve.js
